# os/arch/install
#!/bin/bash

# Before install:
# * Create partition (use fdisk or cfdisk)
#	ex. fdisk /dev/sda
# * Format partition
#   ex. mkfs.ext4 /dev/sda1
# * Format swap partition
#   ex. mkswap /dev/sda2
# * Mount partitions
#	  ex. mount /dev/sda1 /mnt
#	    mount /dev/vol-group/home /mnt/home
#	    mount /dev/sdb5 /mnt/boot

# After install
# * Install grub
#	ex. vim ./grub-setup
#	    sudo ./grub-setup
# * Setup networking
#	ex. ip link
#	    sudo systemctl enable dhcpcd@interface

read -p "Enter hostname: " hostname

lsblk -f
read -p "Enter swap partition: " swapPart

swapon "$swapPart"

# Make sure we have the latest keyring
pacman --noconfirm -Sy archlinux-keyring

# Install base
pacstrap /mnt base base-devel
genfstab -p /mnt >> /mnt/etc/fstab
echo "$hostname" >/mnt/etc/hostname

# Copy setup file and execute
wget https://raw.github.com/Prince86EKNJ/royal-neko/master/os/arch/setup -O /mnt/root/setup
chmod +x /mnt/root/setup
arch-chroot /mnt /bin/bash /root/setup

# Clean up setup file
rm /mnt/root/setup

echo "Done!"
